# package/archw-tools/systemd/user/aw-screen-lock-off.target
[Unit]
Description=Screen lock OFF state off target
