# package/archw-tools/systemd/user/aw-screen-state-on.target
[Unit]
Description=Screen ON state off target
