# package/archw-tools/systemd/user/aw-screen-state-off.target
[Unit]
Description=Screen OFF state off target
