# package/archw-tools/systemd/user/aw-screen-lock-on.target
[Unit]
Description=Screen lock ON state off target
