# package/i3status-rust/systemd/update-checker.timer
[Unit]
Description=System update checker schedule
Requires=update-checker.service

[Timer]
Unit=update-checker.service
OnCalendar=00/6:00
Persistent=true

[Install]
WantedBy=timers.target
