# package/archw-tools/udev/95-aw-monitor-hotplug.rules
SUBSYSTEM=="drm", ACTION=="change", RUN+="/usr/local/bin/aw-udevmonitorevent"
